# templates/liferay.sh.erb
#!/bin/sh
### BEGIN INIT INFO
# Provides:          liferay
# Required-Start:    $network $syslog
# Required-Stop:     $network $syslog
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Starts and stops Liferay (Tomcat)
# Description:       Starts and stops Liferay's Tomcat as if it was a
#                    stardard Debian service: /etc/init.d/liferay start
### END INIT INFO

APPLICATION_NAME='liferay'
LIFERAY_LOCATION="<%= @liferay_location %>"
TOMCAT_HOME="$LIFERAY_LOCATION/tomcat"

getpid() {
    pid=`pgrep -f "java.*liferay"`
}

start() {
    getpid
    if [ -n "$pid" ]; then
        echo "$APPLICATION_NAME (pid $pid) is already running"
        exit 1
    fi
    
    if [ -e $TOMCAT_HOME/catalina_pid.txt ]; then
      rm -f $TOMCAT_HOME/catalina_pid.txt
    fi

    $TOMCAT_HOME/bin/startup.sh
}

stop() {
    $TOMCAT_HOME/bin/shutdown.sh
}

status() {
    getpid
    if [ -n "$pid" ]; then
      echo "$APPLICATION_NAME (pid $pid) is running..."
    else
      echo "$APPLICATION_NAME is NOT running"
      exit 1
    fi
}

case "$1" in
    start)
        start
        ;;
    stop)
        stop
        ;;
    status)
        status
        ;;
    restart)
        stop
        start
        ;;
    *)
        echo $"Usage: $0 {start|stop|restart|status}"
        exit 1
esac

exit 0
